# Docker/workspace/crontab/shakemap
SHELL=/bin/bash

# Build site
* * * * * shakeweb (cd /var/www && echo "-----START-----" && time /usr/bin/flock --verbose -n /tmp/crontabScriptToUpdateEvents.sh.lock /var/www/crontabScriptToUpdateEvents.sh ; echo "-----END-----") >> /var/log/workspace/crontabScriptToUpdateEvents.sh.log 2>&1
#01 00 * * * laradock mv /tmp/crontabScriptToUpdateEvents.sh.log /tmp/crontabScriptToUpdateEvents.sh.log.yesterday
