# Docker/workspace/crontab/shakemap
SHELL=/bin/bash

# Build site for last events
* * * * * laradock (cd /var/www && echo "-----START-----" && time /usr/bin/flock --verbose -n /tmp/crontabScriptToUpdateEvents.sh.lock /var/www/crontabScriptToUpdateEvents.sh ; echo "-----END-----") >> /tmp/crontabScriptToUpdateEvents.sh.log 2>&1
01 00 * * * laradock mv /tmp/crontabScriptToUpdateEvents.sh.log /tmp/crontabScriptToUpdateEvents.sh.log.yesterday

# Rebuild all site every night
00 01 * * * laradock (cd /var/www && echo "-----START-----" && time /usr/bin/flock --verbose -n /tmp/updateEventList.py.lock /usr/bin/python3 /var/www/updateEventList.py ; echo "-----END-----") >> /tmp/updateEventList.py.log 2>&1
00 12 * * * laradock mv /tmp/updateEventList.py.log /tmp/updateEventList.py.log.yesterday
